# conf/zebra.cron
# cron zebra collect once every 5 minutes
MAILTO=""
*/2 * * * * root /usr/bin/zebra -c > /dev/null 2>&1
